# worker-chef.conf
start on filesystem and net-device-up

exec /opt/worker_setup/run_chef.sh
